# packages/vdr/_vdr/system.d/extract-vdr-config.service
[Unit]
Description=Extract VDR sample configuration
DefaultDependencies=false
After=network-online.target graphical.target
Requires=graphical.target
RequiresMountsFor=/storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/vdr-2.6.1/bin/extract_sample_config.sh
StartLimitInterval=0
