# packages/vdr/_vdr/system.d/switch_kodi_vdr.path
[Unit] 
Wants=switch_kodi_vdr.service

[Path] 
PathChanged=/opt/tmp/switch_kodi_vdr

[Install]
WantedBy=multi-user.target
