# packages/vdr/_vdr/system.d/switch_kodi_vdr.service
[Service]  
Type=simple  
ExecStart=/opt/vdr/bin/switch_kodi_vdr.sh

[Install]
WantedBy=multi-user.target
